# packaging/claude-code-control.service
[Unit]
Description=Claude Code Control - tmux controller backend
After=network.target
Requires=network.target

[Service]
Type=simple
EnvironmentFile=-/etc/claude-code-control/config.env
ExecStart=/opt/claude-code-control/venv/bin/uvicorn backend.app.main:app --host ${HOST:-0.0.0.0} --port ${PORT:-8192}
WorkingDirectory=/opt/claude-code-control
Restart=on-failure
RestartSec=5
User=claude-code-control
Group=claude-code-control

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
