# packaging/prerm
#!/bin/bash
set -e

if [ "$1" = "remove" ] && [ -d /run/systemd/system ]; then
    systemctl stop claude-code-control.service 2>/dev/null || true
    systemctl disable claude-code-control.service 2>/dev/null || true
fi
